# filesystem/system-runtime-files-gnome/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Autologin the interactive user on VT tty1 (gnome rootfs only). A VT session is
# attached by logind to the physical seat (seat0), which owns the GPU
# (/dev/dri) and input, so the graphical session runs here (shown in the qemu
# window). This lives in the gnome-only tree so the plain rootfs stays a serial
# console with no graphical autologin. Reset ExecStart first, then override.
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin user --noclear %I $TERM
