# filesystem/system-runtime-files-gnome/usr/lib/systemd/user/gtk4-demo.service
# GTK4 demo as a throwaway graphical-session test client: something to render on
# top of mutter until gnome-shell exists (mutter alone is a black screen).
# Pulled in by mutter.service (Wants); waits for mutter's Wayland socket before
# connecting. PartOf ties it to the session so it stops when the session ends.
[Unit]
Description=GTK4 demo (graphical session test client)
After=mutter.service
Requisite=mutter.service
PartOf=graphical-session.target

[Service]
Type=simple
Environment=GDK_BACKEND=wayland
Environment=WAYLAND_DISPLAY=wayland-0
# No gnome-session here, so the at-spi accessibility bus is not running; disable
# a11y in the client to avoid it stalling on the missing bus.
Environment=GTK_A11Y=none
# mutter (Type=simple) is "started" once spawned, not once its socket exists;
# wait for the socket before launching the client.
ExecStartPre=/bin/sh -c 'until [ -S "${XDG_RUNTIME_DIR}/wayland-0" ]; do sleep 0.1; done'
ExecStart=/usr/bin/gtk4-demo
Restart=on-failure
RestartSec=1
