# filesystem/system-runtime-files-gnome/usr/lib/systemd/user/mutter.service
# Mutter as the graphical session compositor, run as a systemd --user service.
# Stand-in for gnome-session until the shell lands; the launch path (user
# service + devices via logind) is the same one weston used. mutter talks to
# logind directly (no libseat), taking control of seat0 from the tty1 login
# whose identity /etc/profile.d/graphical-session.sh imports into the user
# manager. It runs the native KMS backend as a full display server (not nested)
# and pins the Wayland socket name so the test client can find it.
[Unit]
Description=Mutter Wayland compositor (graphical session)
# graphical-session.target refuses manual start; the compositor brings it up by
# pulling it in as a dependency. The test client is pulled in via Wants.
Wants=graphical-session.target
Before=graphical-session.target
After=graphical-session-pre.target
Wants=gtk4-demo.service

[Service]
Type=simple
Environment=XDG_SESSION_TYPE=wayland
# Cursor theme/size for Wayland clients that inherit this env (mutter itself
# reads org.gnome.desktop.interface cursor-* from gsettings).
Environment=XCURSOR_THEME=Bibata-Modern-Classic
Environment=XCURSOR_SIZE=16
ExecStart=/usr/bin/mutter --wayland --display-server --wayland-display=wayland-0
Restart=always
RestartSec=1
