# filesystem/system-runtime-files/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Autologin the interactive user on VT tty1. Unlike the serial console (ttyS0),
# a VT session is attached by logind to the physical seat (seat0), which owns
# the GPU (/dev/dri) and input devices -- so uaccess grants the user ACLs on
# them and a graphical session can run here (shown in the qemu window). The
# serial console stays a plain debug login. Reset ExecStart first, then override.
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin user --noclear %I $TERM
