# filesystem/system-runtime-files/etc/systemd/system/serial-getty@ttyS0.service.d/autologin.conf
# Autologin root on the serial console (console=ttyS0) for admin/debugging: it
# is the console multiplexed onto the terminal, so it is where diagnostic
# commands are typed. `agetty --autologin root` runs `/bin/login -f root`, which
# skips authentication (root's password is empty and PAM would otherwise reject
# it) and opens a PAM session. The interactive `user` autologs in on the VT
# tty1 (seat0), where the graphical session runs. Reset ExecStart, then override.
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin root --keep-baud 115200,57600,38400,9600 %I $TERM
